# templates/nftables.service.j2
[Unit]
Description=nftables
Documentation=man:nft(8) http://wiki.nftables.org
Wants=network-pre.target
Before=network-pre.target shutdown.target
Conflicts=shutdown.target
DefaultDependencies=no
OnFailure=failmail@%N.service

[Service]
Type=oneshot
RemainAfterExit=yes
StandardInput=null
ProtectSystem=full
ProtectHome=true
ExecStart=/usr/sbin/nft -f /etc/nftables/nftables.nft
ExecReload=/usr/sbin/nft -f /etc/nftables/nftables.nft
ExecStop=/usr/sbin/nft flush ruleset

[Install]
WantedBy=sysinit.target
